# templates/default/resque-manager.upstart.erb
# /etc/init/resque-manager.conf - manage a set of resques

# This example config should work with Ubuntu 12.04+.  It
# allows you to manage multiple resque instances with
# Upstart, Ubuntu's native service management tool.
#
# See resque.conf for how to manage a single resque instance.
#
# Use "stop resque-manager" to stop all resque instances.
# Use "start resque-manager" to start all instances.
# Use "restart resque-manager" to restart all instances.
# Crazy, right?
#

description "Manage the set of resque processes"

# This starts upon bootup and stops on shutdown
start on runlevel [2345]
stop on runlevel [06]

# Set this to the number of resque processes you want
# to run on this machine
env RESQUE_CONF=<%= @resque_conf %>

pre-start script
  pid_start=0
  for i in `cat $RESQUE_CONF`; do
    app=`echo $i | cut -d , -f 1`
    num_workers=`echo $i | cut -d , -f 2`
    queues=`echo $i | cut -d , -f 3-`
    [ ! -z $queues ] || queues='*' # Default to wildcard queue if undefined
    for j in `seq $pid_start $(($pid_start + $num_workers - 1))`; do
      start resque app=$app index=$j queues=$queues
    done
    $pid_start = $j+1
  done
end script
